# microcode_ctl/microcode.rules
KERNEL=="cpu[0-9]*", ACTION=="add", RUN+="/sbin/modprobe microcode"

